# files/DWS4000/copy/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

hwclock -s

if [ -e /dev/sdb1 ] ; then
   mkdir /tmp/Udisk
   mount -t vfat /dev/sdb1 /tmp/Udisk
   chmod 777 /tmp/Udisk
fi

exit 0
